# systemd/lgpowercontrol-update-check.service
[Unit]
Description=LGPowerControl update check

[Service]
Type=oneshot
ExecStart=/opt/lgpowercontrol/lgpowercontrol-update-check.sh
